# AWS/rollapp.service
[Unit]
Description="Rollapp Service"
After=network.target
[Service]
EnvironmentFile=/home/ubuntu/dymension/aws.env
Type=simple
Restart=always
RestartSec=1
User=ubuntu
ExecStart=/bin/bash /home/ubuntu/dymension/dymension-rdk/scripts/run_rollapp.sh
[Install]
WantedBy=multi-user.target